# deploy/osqa.conf
description "Cozy OSQA instance"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5

script
exec sudo sh /home/osqa/osqa/start.sh
end script
